# container_root/etc/profile.d/hadoop.sh
export JAVA_HOME=/usr
export HADOOP_CONF_DIR=/opt/hadoop/etc/hadoop
export SPARK_HOME=/opt/spark/
export PATH=/opt/hadoop/bin:/opt/hsdoop/sbin/:/opt/spark/bin:$PATH
