# container_root/etc/profile.d/hadoop.sh
export JAVA_HOME=/usr/lib/jvm/java-1.8.0-openjdk-amd64
export HADOOP_CONF_DIR=/opt/hadoop/etc/hadoop
export HBASE_CONF_DIR=/opt/hbase/conf
export SPARK_HOME=/opt/spark/
export PATH=/opt/hadoop/bin:/opt/hadoop/sbin/:/opt/spark/bin:/opt/hbase/bin:/opt/jvmtop/:$PATH
